# deploy/roles/telemt_syn_limit/templates/telemt-syn-limit.service.j2
[Unit]
Description=Rate-limit Telemt SYN attempts in Docker forwarding path
Requires=docker.service
After=docker.service network-online.target
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=TELEMT_SYN_INTERFACE={{ ansible_facts["default_ipv4"]["interface"] }}
Environment=TELEMT_SYN_PORT={{ telemt_syn_port }}
Environment=TELEMT_SYN_RATE={{ telemt_syn_rate }}
Environment=TELEMT_SYN_BURST={{ telemt_syn_burst }}
ExecStart=/usr/local/sbin/telemt-syn-limit start
ExecStop=/usr/local/sbin/telemt-syn-limit stop

[Install]
WantedBy=multi-user.target
